# docker_images/pulseaudio/root/etc/services.d/pulseaudio/run
#!/usr/bin/with-contenv bash

# set umask
UMASK_SET=${UMASK_SET:-022}
umask "$UMASK_SET"

exec \
	s6-setuidgid abc /usr/bin/pulseaudio -vvvv --log-time=1 --file /config/docker-pa-server.pa